# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/alsa-utils-alsactl.list
/lib/systemd/system/alsa-state.service
/usr/share/alsa/init/ca0106
/usr/share/alsa/init/test
/lib/systemd/
/lib/udev/rules.d/90-alsa-restore.rules
/lib/systemd/system/basic.target.wants/alsa-restore.service
/lib/systemd/system/
/usr/share/alsa/init/info
/lib/systemd/system/alsa-restore.service
/usr/share/alsa/init/hda
/usr/share/alsa/init/default
/var/lib/alsa/
/lib/systemd/system/basic.target.wants/
/usr/share/alsa/init/00main
/lib/udev/
/usr/share/alsa/init/
/usr/share/alsa/init/help
/lib/udev/rules.d/
/lib/systemd/system/basic.target.wants/alsa-state.service
/usr/sbin/alsactl
