# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/dbus-1.list
/usr/bin/dbus-send
/usr/libexec/dbus-daemon-launch-helper
/usr/share/dbus-1/session.conf
/usr/share/dbus-1/system.conf
/etc/init.d/dbus-1
/usr/bin/dbus-uuidgen
/lib/systemd/system/dbus.socket
/lib/systemd/system/sockets.target.wants/dbus.socket
/usr/bin/dbus-daemon
/lib/systemd/system/dbus.target.wants/dbus.socket
/etc/dbus-1/
/usr/bin/dbus-cleanup-sockets
/usr/share/dbus-1/system.d/
/var/lib/dbus/
/usr/share/dbus-1/session.d/
/usr/share/dbus-1/
/usr/bin/dbus-run-session
/usr/bin/dbus-launch
/usr/share/dbus-1/system-services/
/usr/bin/dbus-update-activation-environment
/lib/systemd/system/dbus.service
/lib/systemd/system/multi-user.target.wants/dbus.service
/etc/default/volatiles/99_dbus
/usr/share/dbus-1/services/
/etc/dbus-1/session.conf
/etc/dbus-1/system.conf
/lib/systemd/system/dbus.target.wants/
/lib/systemd/system/sockets.target.wants/
/usr/bin/dbus-monitor
/lib/systemd/system/multi-user.target.wants/
