# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/opkg.list
/usr/bin/opkg
/usr/share/opkg/intercept/update-modules
/etc/opkg/opkg.conf
/usr/share/opkg/intercept/depmod
/usr/share/opkg/intercept/
/etc/opkg/
/usr/share/opkg/intercept/ldconfig
/lib/systemd/system/opkg-configure.service
/usr/bin/opkg-check-config
/usr/share/opkg/
