# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/run-postinsts.list
/lib/systemd/system/run-postinsts.service
/etc/init.d/run-postinsts
/usr/sbin/run-postinsts
