# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/systemd-extra-utils.list
/bin/systemd-escape
/lib/systemd/system/systemd-hibernate.service
/bin/systemd-tty-ask-password-agent
/usr/bin/systemd-cat
/usr/bin/systemd-stdio-bridge
/usr/bin/systemd-delta
/lib/systemd/system/systemd-initctl.socket
/lib/systemd/system/systemd-suspend.service
/lib/systemd/system/sleep.target
/lib/systemd/system-sleep/
/lib/systemd/system/systemd-ask-password-console.service
/lib/systemd/system/systemd-ask-password-wall.path
/lib/systemd/system/systemd-ask-password-wall.service
/lib/systemd/systemd-initctl
/bin/systemd-ask-password
/lib/systemd/system/systemd-hybrid-sleep.service
/lib/systemd/systemd-reply-password
/lib/systemd/system/sockets.target.wants/systemd-initctl.socket
/lib/systemd/system/systemd-nspawn@.service
/usr/bin/systemd-cgtop
/lib/systemd/systemd-socket-proxyd
/usr/bin/systemd-nspawn
/usr/bin/systemd-detect-virt
/usr/bin/systemd-run
/lib/systemd/system/systemd-initctl.service
/lib/systemd/system/multi-user.target.wants/systemd-ask-password-wall.path
/lib/systemd/systemd-ac-power
/bin/systemd-inhibit
/lib/systemd/system-generators/systemd-gpt-auto-generator
/usr/lib/tmpfiles.d/systemd-nspawn.conf
/lib/systemd/systemd-sleep
/usr/bin/systemd-cgls
/lib/systemd/systemd-cgroups-agent
/lib/systemd/system/sysinit.target.wants/systemd-ask-password-console.path
/lib/systemd/system/systemd-ask-password-console.path
/usr/bin/systemd-path
