# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/systemd-vconsole-setup.list
/lib/systemd/systemd-vconsole-setup
/lib/systemd/system/sysinit.target.wants/systemd-vconsole-setup.service
/lib/systemd/system/systemd-vconsole-setup.service
