# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/ti-ipc.list
/usr/lib/libtitransportrpmsg.so.1.0.0
/usr/bin/fault
/usr/bin/MessageQMultiMulti
/usr/lib/libtiipcutils.so.1
/usr/bin/lad_dra7xx
/usr/bin/Msgq100
/usr/lib/libtiipcutils_lad.so.1
/usr/bin/MessageQBench
/etc/init.d/tiipclad-daemon.sh
/usr/lib/libtiipcutils.so.1.0.0
/usr/lib/libmmrpc.so.1.0.0
/usr/lib/libtitransportrpmsg.so.1
/usr/lib/libtiipcutils_lad.so.1.0.0
/usr/lib/libtiipc.so.1
/usr/lib/libmmrpc.so.1
/usr/lib/libtiipc.so.1.0.0
