# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/ti-ipc.preinst
#!/bin/sh
if [ -z "$D" -a -f "/etc/init.d/tiipclad-daemon.sh" ]; then
	/etc/init.d/tiipclad-daemon.sh stop
fi
if type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-f -r $D"
	else
		OPT="-f"
	fi
	update-rc.d $OPT tiipclad-daemon.sh remove
fi
