# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/udev.list
/etc/init.d/systemd-udevd
/lib/udev/rules.d/60-persistent-input.rules
/lib/udev/rules.d/70-uaccess.rules
/lib/udev/rules.d/60-drm.rules
/usr/share/bash-completion/completions/
/etc/udev/scripts/usb2-rules.sh
/lib/udev/rules.d/60-block.rules
/lib/udev/rules.d/99-systemd.rules
/lib/udev/rules.d/73-seat-late.rules
/lib/udev/scsi_id
/lib/udev/rules.d/60-persistent-v4l.rules
/lib/systemd/system/initrd-udevadm-cleanup-db.service
/lib/udev/cdrom_id
/lib/udev/rules.d/60-persistent-storage.rules
/lib/udev/rules.d/50-udev-default.rules
/etc/udev/rules.d/local.rules
/lib/systemd/system/systemd-udevd.service
/etc/udev/hwdb.d/
/etc/udev/rules.d/touchscreen.rules
/lib/udev/rules.d/60-cdrom_id.rules
/lib/systemd/system/systemd-udevd-kernel.socket
/lib/udev/rules.d/90-vconsole.rules
/lib/udev/rules.d/60-persistent-storage-tape.rules
/lib/udev/rules.d/80-net-setup-link.rules
/etc/udev/
/lib/systemd/system/systemd-udev-trigger.service
/lib/systemd/system/sysinit.target.wants/systemd-udevd.service
/sbin/udevd
/lib/udev/rules.d/60-serial.rules
/usr/share/bash-completion/completions/udevadm
/lib/systemd/system/systemd-udevd-control.socket
/lib/udev/rules.d/60-persistent-alsa.rules
/lib/udev/mtd_probe
/lib/udev/ata_id
/lib/udev/collect
/lib/udev/rules.d/60-evdev.rules
/lib/systemd/system/sysinit.target.wants/systemd-udev-trigger.service
/lib/udev/rules.d/70-power-switch.rules
/usr/share/bash-completion/
/lib/udev/rules.d/70-mouse.rules
/lib/udev/rules.d/50-firmware.rules
/etc/udev/scripts/
/lib/udev/rules.d/78-sound-card.rules
/lib/udev/rules.d/64-btrfs.rules
/bin/udevadm
/etc/udev/udev.conf
/etc/udev/scripts/usb1-rules.sh
/etc/udev/rules.d/
/lib/systemd/system/systemd-udev-settle.service
/lib/systemd/system/sockets.target.wants/systemd-udevd-kernel.socket
/lib/udev/rules.d/75-probe_mtd.rules
/lib/udev/v4l_id
/lib/systemd/systemd-udevd
/lib/udev/rules.d/75-net-description.rules
/lib/systemd/system/sockets.target.wants/systemd-udevd-control.socket
/lib/udev/rules.d/71-seat.rules
/lib/udev/rules.d/80-drivers.rules
