# AM5728/ti-processor-sdk-linux-am57xx-evm-03.02.00.05/linux-devkit/sysroots/armv7ahf-neon-linux-gnueabi/var/lib/opkg/info/volatile-binds.list
/sbin/mount-copybind
/lib/systemd/system/var-volatile-lib.service
